# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:backend]
command=python -m uvicorn main:app --host 127.0.0.1 --port 5500
directory=/app/backend
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=PYTHONPATH="/app/backend",OPENROUTER_API_KEY="%(ENV_OPENROUTER_API_KEY)s",OPENROUTER_MODEL="%(ENV_OPENROUTER_MODEL)s"

[program:frontend]
command=/app/frontend/node_modules/.bin/next start -p 5001
directory=/app/frontend
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
